# docker/supervisor/reportflow.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -F
priority=10
autorestart=true
autostart=true

[program:nginx]
command=nginx -g "daemon off;"
priority=20
autorestart=true
autostart=true

[program:queue]
command=php /var/www/html/artisan queue:work --tries=3 --timeout=120
autorestart=true
autostart=true

[program:scheduler]
command=/bin/sh -c "while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done"
autorestart=true
autostart=true